# build-tests/x86/rawhide/test-image-nitro-enclave/root/usr/lib/systemd/system/vsock_server.service
[Unit]
Description=Vsock Server for Nitro Enclaves
After=network-online.target nitro-enclave-alive.service

[Service]
ExecStart=/usr/bin/vsock_server
Restart=on-failure

[Install]
WantedBy=multi-user.target
