# build-tests/x86/tumbleweed/test-image-agama/root/etc/systemd/system/agama-certificate-issue.service
[Unit]
Description=Generate issue file for Agama SSL certificate

[Service]
Type=oneshot
ExecStart=agama-issue-generator --ssl

[Install]
WantedBy=default.target
