# build-tests/x86/tumbleweed/test-image-agama/root/etc/systemd/system/agama-avahi-issue.service
[Unit]
Description=Generate issue file for Agama URL from Avahi

After=avahi-daemon.service

[Service]
ExecStart=agama-issue-generator --watch-avahi
Type=simple

[Install]
WantedBy=default.target
