# build-tests/x86/tumbleweed/test-image-agama/root/etc/systemd/system/agama-certificate-issue.path
[Unit]
Description=Watch the Agama SSL certificate for changes
Before=systemd-user-sessions.service

[Path]
Unit=agama-certificate-issue.service
PathChanged=/etc/agama.d/ssl/cert.pem

[Install]
WantedBy=default.target
